# deploy/systemd/trueid-backup.service
[Unit]
Description=TrueID Database Backup

[Service]
Type=oneshot
User=trueid
WorkingDirectory=/opt/trueid
ExecStart=/opt/trueid/scripts/backup.sh /opt/trueid/backups
